# debian/python3-keylime-agent.keylime_agent.path
[Unit]
Description=Monitor changes in the config for agent

[Path]
PathExists=/etc/keylime.conf
Unit=keylime_agent.service

[Install]
WantedBy=multi-user.target
